# deploy/sunyun-portal.service
[Unit]
Description=Sunyun Portal
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/sunyun-portal
EnvironmentFile=/opt/sunyun-portal/.env
ExecStart=/usr/bin/node /opt/sunyun-portal/server.js
Restart=always
RestartSec=5
User=www-data
Group=www-data

[Install]
WantedBy=multi-user.target
